# nginx-php-fpm-sshd/conf/crontab.cron
* * * * * sh /var/run/check_crontab.sh >> /var/log/check_crontab.log
